# nocout_device_app/poller/ETL/cron.d/historical_clean_data
0 23 * * * /usr/bin/python nocout/mysql_clean.py -a historical_data
