# nocout_device_app/poller/ETL/cron.d/live_clean_data
0 23 * * * /usr/bin/python nocout/mysql_clean.py -a live_data 
