# nocout_device_app/pollers_ETL/cron.d/historical_data
#---------------------------------------------------------------------------------------------------
## Live Mongodb --> Historical Mongodb

# Bihourly scripts
*/30 * * * * /usr/bin/python nocout/performance/service/service_mongo_aggregation_half_hourly.py -t 30 -s network_perf -d performance_performancenetworkbihourly
*/30 * * * * /usr/bin/python nocout/performance/service/service_mongo_aggregation_half_hourly.py -t 30 -s service_perf -d performance_performanceneservicebihourly

# Hourly scripts
0 * * * * /usr/bin/python nocout/performance/service/aggregation_all.py -t 1 -f hourly -s performance_performancenetworkbihourly -d performance_performancenetworkhourly
0 * * * * /usr/bin/python nocout/performance/service/mongo_aggregation_all.py -t 1 -f hourly -s performance_performanceservicebihourly -d performance_performanceservicehourly

# Daily scripts
0 23 * * * /usr/bin/python nocout/performance/service/aggregation_all.py -t 24 -f daily -s prformance_performancenetworkhourly -d performance_performancenetworkdaily
0 23 * * * /usr/bin/python nocout/performance/service/aggregation_all.py -t 24 -f daily -s performance_performanceservicehourly -d performance_performanceservicedaily
0 23 * * * /usr/bin/python nocout/performance/interface/mongo_aggregation_freq_based_base.py -t 24 -f daily -s status_perf -d interface_perf_daily
0 23 * * * /usr/bin/python nocout/performance/inventory/inventory_mongo_aggregation_daily.py

# Weekly scripts
0 23 * * 0 /usr/bin/python nocout/performance/service/aggregation_all.py -t 168 -f weekly -s performance_performancenetworkdaily -d performance_performancenetworkweekly
0 23 * * 0 /usr/bin/python nocout/performance/service/aggregation_all.py -t 168 -f weekly -s performance_performanceservicedaily -d performance_performanceserviceweekly
0 23 * * 0 /usr/bin/python nocout/performance/service/aggregation_all.py -t 168 -f weekly -s performance_performancestatusdaily -d performance_performancestatusweekly
0 23 * * 0 /usr/bin/python nocout/performance/interface/mongo_aggregation_freq_based_base.py -t 168 -f weekly -s nocout_inventory_service_perf_data -d inventory_perf_weekly


## Historical Mongodb --> Historical Mysql

# Daily Migration
5 23 * * * /usr/bin/python nocout/performance/service/historical_mongo_to_mysql.py -t 24 -s interface_perf_daily -d performance_performancestatusdaily

# Weekly Migration
5 23 * * 0 /usr/bin/python nocout/performance/service/historical_mongo_to_mysql.py -t 168 -s inventory_perf_weekly -d performance_performanceinventoryweekly
#-------------------------------------------------------------------------------------------------------
